# systemd/logline-dualzone-logger.service.example
[Unit]
Description=Logline dual-zone greenhouse logger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=logline
WorkingDirectory=/opt/logline-greenhouse-edge
EnvironmentFile=-/opt/logline-greenhouse-edge/.env
ExecStart=/opt/logline-greenhouse-edge/scripts/run_dualzone_logger.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
