# systemd/logline-greenhouse-edge.service.example
[Unit]
Description=Logline greenhouse edge bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=logline
WorkingDirectory=/path/to/logline-greenhouse-edge
EnvironmentFile=-/path/to/logline-greenhouse-edge/.env
ExecStart=/path/to/logline-greenhouse-edge/scripts/run_edge_bridge.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
